# .docker/build/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
